# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

# ── bgutil Proof-of-Origin Token server ───────────────────────────────────────
# Runs on port 4416 (localhost only). yt-dlp plugin picks up tokens
# automatically — no extra config needed in bot.py.
# This fixes the "Sign in to confirm you're not a bot" error on cloud IPs.
[program:bgutil-pot]
command=node /opt/bgutil/server/build/main.js
directory=/opt/bgutil/server
autorestart=true
startretries=5
priority=10
stdout_logfile=/var/log/supervisor/bgutil.log
stderr_logfile=/var/log/supervisor/bgutil.log

# ── Telegram bot ───────────────────────────────────────────────────────────────
[program:bot]
command=python -u /app/bot.py
autorestart=true
priority=20
startsecs=5
stdout_logfile=/var/log/supervisor/bot.log
stderr_logfile=/var/log/supervisor/bot.log
